# roles/auto-ebs/templates/usr/lib/systemd/system/auto.ebs-detach.path
# $Id$
[Unit]
Description=PathModified=/ebs/ /etc/auto.ebs-detach.sh

[Path]
PathModified=/ebs/

[Install]
WantedBy=multi-user.target
