# roles/auto.ebs/templates/usr/lib/systemd/system/auto.ebs-detach.service
# $Id$
[Unit]
Description=/etc/auto.ebs-detach.sh
After=autofs.service

[Service]
Type=oneshot
ExecStart=/etc/auto.ebs-detach.sh

[Install]
WantedBy=multi-user.target
